# config/systemd/user/ddns.service
[Unit]
Description=Update DDNS Info

[Service]
Type=oneshot
ExecStart=/home/remilia/.config/systemd/user/ddns.py
