# config/systemd/user/ddns.timer
[Unit]
Description=Update DDNS Info

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
